# userspace_driver/elan1200.service
[Unit]
Description=ELAN1200 userspace driver

[Service]
Type=simple
ExecStart=/usr/local/bin/hid_elan1200

[Install]
WantedBy=multi-user.target
